# deploy/systemd/maintenance.timer
[Unit]
Description=Run Machine Image Uploader maintenance periodically

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
